# confluent-kafka.service.j2
[Unit]
Description=Confluent Kafka Broker (Confluent {{ confluent_package_version }})
Documentation=https://docs.confluent.io
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ broker_user }}
Group={{ broker_group }}

# Source environment file (JVM, JMX, Prometheus opts set here)
EnvironmentFile=-{{ broker_install_dir }}/bin/broker-env.sh

# Explicit environment overrides (mirrors zookeeper.service.j2 pattern)
Environment="KAFKA_HEAP_OPTS={{ broker_heap_opts }}"
Environment="LOG_DIR={{ broker_log_dir }}"
{% if prometheus_jmx_exporter_enabled %}
Environment="KAFKA_OPTS={{ java_agent_string }}"
{% endif %}

# Start Command
ExecStart={{ broker_install_dir }}/bin/kafka-server-start {{ broker_config_dir }}/server.properties

# Stop Command
ExecStop={{ broker_install_dir }}/bin/kafka-server-stop

# Restart Policy
Restart=on-failure
RestartSec=10
TimeoutStopSec=180

# Resource Limits
LimitNOFILE=100000
LimitNPROC=32768

# Working Directory
WorkingDirectory={{ broker_install_dir }}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=confluent-kafka

[Install]
WantedBy=multi-user.target
